# rw/bind-dirs/etc/systemd/system/openvpn@.service.d/override.conf
[Service]
Group=openvpn
WorkingDirectory=/etc/openvpn/client
